# 99-gamepad_watchdog.rules
SUBSYSTEM=="input", KERNEL=="event[0-9]*", RUN+="/usr/bin/python /userdata/RetroFlag/GamepadWatchdog.py refresh"
